# kritomatic_xremap/xremap-kritomatic.path
[Unit]
Description=Watch Kritomatic Directory for Changes

[Path]
PathChanged=REPODIR_PLACEHOLDER/
PathChanged=%h/xremap_config.yaml
Unit=xremap-kritomatic-restart.service

[Install]
WantedBy=default.target
