# kritomatic_xremap/xremap-kritomatic.path
[Unit]
Description=Watch Kritomatic and Personal Config for changes

[Path]
# Watch the directory itself to catch new file writes/moves
PathModified=REPODIR_PLACEHOLDER/
# Also watch your personal overrides
PathModified=%h/xremap_config.yaml
Unit=xremap-kritomatic.service

[Install]
WantedBy=graphical-session.target
